# scripts/com.esp32-round-clock.usb-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.esp32-round-clock.usb-daemon</string>
  <key>ProgramArguments</key>
  <array>
    <string>HOME_PLACEHOLDER/Library/Application Support/esp32-round-clock/.venv/bin/python3</string>
    <string>HOME_PLACEHOLDER/Library/Application Support/esp32-round-clock/usb_daemon.py</string>
    <string>--run</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>HOME_PLACEHOLDER/Library/Application Support/esp32-round-clock/usb-daemon.log</string>
  <key>StandardErrorPath</key>
  <string>HOME_PLACEHOLDER/Library/Application Support/esp32-round-clock/usb-daemon.log</string>
</dict>
</plist>
